# supervisord.conf
[supervisord]
nodaemon=true

[program:main]
command=/usr/local/bin/openclaw gateway --port %(ENV_OPENCLAW_GATEWAY_PORT)s
user=openclaw
environment=HOME="/home/openclaw",USER="openclaw"
priority=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/openclaw/openclaw.log        ; 指定标准输出日志路径
stdout_logfile_maxbytes=50MB                      ; 单个日志文件最大50MB（可选）
stdout_logfile_backups=10                          ; 保留10个备份（可选）
stdout_capture_maxbytes=1MB                        ; 捕获模式最大字节（可选）
stderr_logfile=/var/log/openclaw/openclaw.err.log     ; 错误日志路径（可选）
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10

[program:browser]
command=/opt/conda/bin/python /scripts/browser_daemon/browser_daemon.py
user=openclaw
environment=HOME="/home/openclaw",USER="openclaw"
priority=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/openclaw/browser_daemon.log        ; 指定标准输出日志路径
stdout_logfile_maxbytes=50MB                      ; 单个日志文件最大50MB（可选）
stdout_logfile_backups=10                          ; 保留10个备份（可选）
stdout_capture_maxbytes=1MB                        ; 捕获模式最大字节（可选）
stderr_logfile=/var/log/openclaw/browser_daemon.err.log     ; 错误日志路径（可选）
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10